# package/postinstall
#!/bin/bash

/etc/init.d/myapp start
chkconfig myapp on
